# etc/cron.daily/disable-crashpad.sh
#!/bin/sh

chmod 644 /opt/brave/chrome_crashpad_handler
chmod 644 /opt/slack/chrome_crashpad_handler
chmod 644 /opt/vscodium/chrome_crashpad_handler
chmod 644 /opt/google/chrome/chrome_crashpad_handler
